# docker-nginx/apps.conf
[program:uwsgi]
command = /usr/local/bin/uwsgi --ini /app/uwsgi.ini
user=admin

[program:nginx]
command = /usr/sbin/nginx
